# deploy/ndld-update.service
[Unit]
Description=Update ndld Docker container
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/ndld/update.sh
Environment=COMPOSE_DIR=/opt/ndld
Environment=HEALTH_URL=http://localhost:8080/health

# Security hardening
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/opt/ndld

[Install]
WantedBy=multi-user.target
